# crontab
0 0 * * * python /home/scrapy/geoscrap/main.py > /proc/1/fd/1 2>&1
